# rpcd/rpcd.service
[Unit]
Description=OpenWrt ubus RPC daemon

[Service]
ExecStart=/usr/bin/rpcd -t 30

[Install]
WantedBy=multi-user.target
